# mtp-responder.service
[Unit]
Description=MTP responder

[Service]
Type=simple
ExecStart=/usr/bin/mtp-responder
KillMode=process